# dockerize/production/systemd/goldentempo-backup.service
[Unit]
Description=Golden Tempo Postgres backup (pg_dump -> gzip -> rclone off-site)
Documentation=file:///opt/goldentempo/restore.md
# Only run once the Docker daemon is up — backup.sh talks to the postgres
# container via docker compose exec.
Requires=docker.service
After=docker.service network-online.target

[Service]
Type=oneshot
# EnvironmentFile carries the backup knobs (BACKUP_DIR, RETENTION_DAYS,
# RCLONE_REMOTE, BACKUP_HEARTBEAT_FILE, ...) alongside the app config.
EnvironmentFile=/opt/goldentempo/.env
ExecStart=/opt/goldentempo/backup.sh
# A backup should never run forever; fail loudly instead.
TimeoutStartSec=1800
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
